# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/polymarket/supervisord.log
pidfile=/var/run/supervisord.pid
user=polymarket

[program:scanner]
command=python -m agents.scanner
directory=/opt/polymarket-system
autostart=true
autorestart=true
startretries=10
stderr_logfile=/var/log/polymarket/scanner.err.log
stdout_logfile=/var/log/polymarket/scanner.log
environment=PYTHONPATH="/opt/polymarket-system"

[program:analyst]
command=python -m agents.analyst
directory=/opt/polymarket-system
autostart=true
autorestart=true
startretries=10
stderr_logfile=/var/log/polymarket/analyst.err.log
stdout_logfile=/var/log/polymarket/analyst.log
environment=PYTHONPATH="/opt/polymarket-system"

[program:lp_manager]
command=python -m agents.lp_manager
directory=/opt/polymarket-system
autostart=true
autorestart=true
startretries=10
stderr_logfile=/var/log/polymarket/lp_manager.err.log
stdout_logfile=/var/log/polymarket/lp_manager.log
environment=PYTHONPATH="/opt/polymarket-system"

[program:notifier]
command=python -m agents.notifier
directory=/opt/polymarket-system
autostart=true
autorestart=true
startretries=10
stderr_logfile=/var/log/polymarket/notifier.err.log
stdout_logfile=/var/log/polymarket/notifier.log
environment=PYTHONPATH="/opt/polymarket-system"

[program:executor]
command=python -m agents.executor
directory=/opt/polymarket-system
autostart=true
autorestart=true
startretries=10
stderr_logfile=/var/log/polymarket/executor.err.log
stdout_logfile=/var/log/polymarket/executor.log
environment=PYTHONPATH="/opt/polymarket-system"

[program:daily_report]
command=python scripts/daily_report.py
directory=/opt/polymarket-system
autostart=false
autorestart=false
; Run via cron: 0 0 * * * supervisorctl start daily_report
stderr_logfile=/var/log/polymarket/daily_report.err.log
stdout_logfile=/var/log/polymarket/daily_report.log
environment=PYTHONPATH="/opt/polymarket-system"
